# data/export/upstart/process.conf.erb
start on starting <%= app %>-<%= name %>
stop on stopping <%= app %>-<%= name %>
respawn

env PORT=<%= port %>
<% engine.env.each do |name,value| -%>
env <%= name.upcase %>='<%= value.gsub(/'/, "'\"'\"'") %>'
<% end -%>

setuid <%= user %>

chdir <%= engine.root %>

script
  <%= engine.script_preamble if engine.respond_to?(:script_preamble) %>
  <%= process.command %>
  <%= engine.script_postamble if engine.respond_to?(:script_postamble) %>
end script
